# dotdrop/dotfiles/etc/systemd/system/auto-hibernate.timer
[Unit]
Description=Auto hibernate system when the battery is low

[Timer]
OnUnitActiveSec=1min
OnBootSec=10s

[Install]
WantedBy=timers.target
