# dotfiles/etc/systemd/system/auto-hibernate.service
[Unit]
Description=Hibernate system when the battery is low

[Service]
Type=oneshot
ExecStart=/bin/bash /bin/hibernate-if-low.sh
