# dotfiles/etc/systemd/system/monitor-battery.service
[Unit]
Description=Monitors the battery level and responds accordingly

[Service]
Type=oneshot
ExecStart=/bin/bash /bin/monitor-battery.sh
