# dotfiles/etc/systemd/system/monitor-battery.timer
[Unit]
Description=Monitors the battery level and responds accordingly

[Timer]
OnUnitActiveSec=1min
OnBootSec=10s

[Install]
WantedBy=timers.target
